# utility/port_test_deamon/port_test_daemon.service
[Unit]
Description=A daemon to test a open port

[Service]
ExecStart=/usr/bin/port_test.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
